# manifests/website.service
[Unit]
Description=BusyBox httpd Server

[Service]
Type=simple
user=root
ExecStart=/usr/bin/WebServer.sh

[Install]
WantedBy=multi-user.target